# etc/services.d/crond/run
#!/usr/bin/execlineb -P
/usr/sbin/cron -f -l 2
